# rootfs/etc/s6-overlay/s6-rc.d/backup_api/run
#!/usr/bin/with-contenv bashio

bashio::log.info "---------Starting Backup script--------"

declare FleetAssistantServerIP
declare FleetToken
declare Installation_id
declare backup_enable

FleetAssistantServerIP="10.170.204.1"

# Get Fleet Assistant backup enabled status from addon configuration page
if bashio::config.has_value "server.enable_backup"; then
    backup_enable=$(bashio::config "server.enable_backup")
fi

if [[ $backup_enable = false ]]; then
    bashio::log.warning "Backup is disabled. Backup creation will not happen."
    bashio::log.warning "Backup function will stop. Restart addon with backup option enabled in config to start backup functionality again"
    while true; do sleep 1h; done
fi

# Get Fleet Assistant token from addon configuration page
if bashio::config.has_value "server.fleet_token"; then
    FleetToken=$(bashio::config "server.fleet_token")
fi
# Get Fleet Assistant Installation ID from addon configuration page
if bashio::config.has_value "server.installation_id"; then
    Installation_id=$(bashio::config "server.installation_id")
fi



exec python3 -u /etc/s6-overlay/scripts/check_backup_server_status.py --FleetAssistantServerIP ${FleetAssistantServerIP} --FleetToken ${FleetToken} --Installation_id ${Installation_id}